# container/supervisord.conf
[supervisord]
nodaemon=true

[program:bro]
command=/usr/local/bro/bin/broctl deploy