# container/supervisord.conf
[supervisord]
nodaemon=true

[program:bro]
command=/bin/bash -c "sed -i 's/${INTERFACE}/'$INTERFACE' /g' /usr/local/bro/etc/node.cfg && /usr/local/bro/bin/broctl deploy"