# container/supervisord.conf
[supervisord]
nodaemon=true

[program:zeek]
command=/bin/bash -c "sed -i 's/${INTERFACE}/'$INTERFACE' /g' /usr/local/zeek/etc/node.cfg && /usr/local/zeek/bin/zeekctl deploy"
